# templates/init/systemd/kong.service.erb
# This file is managed by Puppet
[Unit]
Description=Kong daemon

[Service]
<%- if !@nginx_conf.nil? -%>
ExecStart=<%= @kong_path -%> start --conf @config_file_path --nginx-conf @nginx_conf
ExecReload=<% @kong_path -%> reload --conf @config_file_path --nginx-conf @nginx_conf
<% else -%>
ExecStart=<%= @kong_path -%> start --conf @config_file_path
ExecReload=<% @kong_path -%> reload --conf @config_file_path --nginx-conf
<%- end -%>
ExecStop<%= @kong_path -%> stop
Restart=on-failure