# templates/init/systemd/kong.service.erb
# This file is managed by Puppet
[Unit]
Description=Kong daemon

[Service]
ExecStart=<%= @kong_path -%> start
ExecReload=<%= @kong_path -%> reload
ExecStop=<%= @kong_path -%> stop
Restart=on-failure
