# templates/init/upstart/kong.conf.erb
# This file was generated by Puppet

description	"Kong daemon"

start on runlevel [2345]
stop on runlevel [!2345]

script
	exec <%= @kong_path %> start
end script
